# 7.1.0/crontab

PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * /var/www/html/cron/vtigercron.sh > /var/www/html/logs/vtgiercrm_cron.log

