# 5.2.0/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
* * * * * /var/www/html/vtiger/cron/vtigercron.sh &> /var/www/html/vtiger/logs/vtgier_cron.log
