# versions/apache/php7/debian/7.1.0/crontab
* * * * * root docker-vtiger-cron.sh timestamp
* * * * * root docker-vtiger-cron.sh scheduler
* * * * * root docker-vtiger-cron.sh localhost-proxy
